# root/etc/s6-overlay/s6-rc.d/svc-tigervncserver/run
#!/usr/bin/with-contenv bash

#exec s6-setuidgid user \
#  tigervncserver -xstartup /usr/bin/plasma_session

runuser -l user -c \
'tigervncserver -xstartup /usr/bin/plasma_session'
