# root/etc/s6-overlay/s6-rc.d/init-tigervncserver/run
#!/usr/bin/with-contenv bash

#exec s6-setuidgid user \
#  printf "Pa55w.rd\nPa55w.rd\n\n" | vncpasswd

runuser -l user -c \
'printf "Pa55w.rd\nPa55w.rd\n\n" | vncpasswd'
